# my_work/part_2/opt/fleetsie/custom_post/leds-off.service
[Unit]
Description=turn off Pi leds completely
[Service]
Type=oneshot
ExecStart=/opt/ampi/leds_off
[Install]
WantedBy=multi-user.target
